# smtp-relay.service
[Unit]
Description=SMTP Relay Service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/smtp-relay
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
TimeoutStartSec=0
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 